# deploy/systemd/neko-what-if-refresh.timer
[Unit]
Description=Hourly refresh of Neko's offline What If event cache

[Timer]
OnBootSec=2min
OnUnitActiveSec=1h
AccuracySec=1min
Persistent=true
Unit=neko-what-if-refresh.service

[Install]
WantedBy=timers.target
