# deploy/systemd/neko-tts.service
[Unit]
Description=Neko resident English KittenTTS worker
After=local-fs.target
Wants=local-fs.target
StartLimitIntervalSec=120s
StartLimitBurst=3

[Service]
Type=notify
NotifyAccess=main
User=neko
Group=neko
WorkingDirectory=/home/neko/repos/nekokitty
RuntimeDirectory=neko
RuntimeDirectoryMode=0750
Environment=PULSE_SERVER=unix:/dev/null
ExecStart=/home/neko/.local/share/neko/venvs/kittentts/bin/python /home/neko/repos/nekokitty/scripts/serve_kittentts.py --socket /run/neko/tts.sock --profile micro --threads 4
Restart=on-failure
RestartSec=5s
TimeoutStartSec=45s
TimeoutStopSec=15s
NoNewPrivileges=true
PrivateTmp=true
PrivateDevices=true
PrivateNetwork=true
ProtectSystem=strict
ProtectHome=read-only
ProtectClock=true
ProtectControlGroups=true
ProtectHostname=true
ProtectKernelLogs=true
ProtectKernelModules=true
ProtectKernelTunables=true
ProtectProc=invisible
ProcSubset=pid
RestrictAddressFamilies=AF_UNIX
RestrictNamespaces=true
RestrictRealtime=true
RemoveIPC=true
SystemCallArchitectures=native
LockPersonality=true
MemoryDenyWriteExecute=false
RestrictSUIDSGID=true
CapabilityBoundingSet=
AmbientCapabilities=
MemoryHigh=400M
MemoryMax=512M
MemorySwapMax=0
OOMScoreAdjust=400
CPUWeight=35
Nice=5
TasksMax=64
UMask=0027

[Install]
WantedBy=multi-user.target
