# deploy/systemd/bluetooth.service.d/zz-neko-audio.conf
[Service]
# NVIDIA's R39 override disables the BlueZ audio, A2DP, and AVRCP plugins.
# Keep the Serial Port Profile disabled, but restore Bluetooth speaker audio.
ExecStart=
ExecStart=/usr/libexec/bluetooth/bluetoothd -d --noplugin=sap
