# deploy/systemd/user/neko-voice-assistant.service
[Unit]
Description=Neko attended voice assistant
After=neko-audio-policy.service network.target
Wants=neko-audio-policy.service

[Service]
Type=simple
WorkingDirectory=%h/repos/nekokitty
Environment=PYTHONUNBUFFERED=1
ExecStart=%h/.local/share/neko/venvs/asr/bin/python %h/repos/nekokitty/scripts/neko_voice_assistant.py --device pipewire --language en --attended-cat-sounds --verbose-transcripts
Restart=on-failure
RestartSec=3
TimeoutStopSec=15
KillMode=mixed

[Install]
WantedBy=default.target
