# deploy/systemd/user/neko-audio-policy.service
[Unit]
Description=Neko ReSpeaker-first audio routing policy
After=pipewire.service pipewire-pulse.service wireplumber.service
Wants=pipewire.service pipewire-pulse.service wireplumber.service

[Service]
Type=notify
NotifyAccess=main
WorkingDirectory=%h/repos/nekokitty
Environment=PYTHONUNBUFFERED=1
ExecStart=/usr/bin/python3 %h/repos/nekokitty/scripts/neko_audio_policy.py
Restart=on-failure
RestartSec=3
TimeoutStartSec=30
TimeoutStopSec=10
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=read-only
RestrictAddressFamilies=AF_UNIX
LockPersonality=true
RestrictSUIDSGID=true
TasksMax=32
UMask=0077

[Install]
WantedBy=default.target
